# gentoo/runscript
#!/sbin/runscript

extra_commands="checkconfig"
PIDFILE=/var/run/twod.pid
CONFFILE=/etc/twod/twodrc

depend() {
    need net
}

checkconfig() {
    if ! [ -f ${CONFFILE} ]; then
        eerror "You need to setup ${CONFFILE} to run twod"
        eerror "An example can be found in /usr/share/doc/twod"
        return 1
    fi

    return 0
}

start() {
    checkconfig || return 1

    ebegin "Starting twod"
    start-stop-daemon --start --exec /usr/bin/twod --pidfile "${PIDFILE}"
    eend $? "Failed to start twod"
}

stop() {
    ebegin "Stopping twod"
    start-stop-daemon --stop --pidfile "${PIDFILE}" --exec /usr/bin/twod
    eend $? "Failed to stop twod"
}
